# crontab
0 8 * * 2 /usr/local/bin/scrapper > /proc/1/fd/1 2>/proc/1/fd/2
5 8 * * 2 /usr/local/bin/select_movie > /proc/1/fd/1 2>/proc/1/fd/2
0 8 * * 1 /usr/local/bin/select_day > /proc/1/fd/1 2>/proc/1/fd/2
